# hangman/conf/gunicorn.conf
# /etc/init/gunicorn.conf
description "Gunicorn application server handling hangman"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid www-data
chdir /home/ubuntu/games/hangman

exec /usr/local/bin/gunicorn --workers 3 --bind unix:/home/ubuntu/games/hangman/hangman.sock hangman.wsgi:application
